# ClusterPod/crontab/crontab
0 23    * * *   ubuntu sh /home/ubuntu/docker-reg-cred/recreate-docker-registry.sh 2>&1 /home/ubuntu/docker-reg-cred/output.log